# s6/services.d/twitter_comic_streamer/run
#!/usr/bin/with-contenv sh
/app/twitter_comic_streamer
